# systemd/system/tmux-save.service
[Unit]
Description=Save tmux sessions before shutdown/sleep
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target sleep.target

[Service]
Type=oneshot
User=batman
Environment=HOME=/home/batman
ExecStart=/bin/bash -c 'tmux -S /tmp/tmux-1000/default list-sessions &>/dev/null && /home/batman/.config/tmux/plugins/tmux-resurrect/scripts/save.sh'
RemainAfterExit=yes
TimeoutStartSec=15

[Install]
WantedBy=shutdown.target reboot.target halt.target sleep.target
